# Files/rc.local
#!/bin/sh -e
docker start apt-cacher-ng-run

exit 0